# deploy/agent.service
[Unit]
Description=__AGENT_DISPLAY_NAME__ — Telegram ↔ Claude agent (long-poll bot + persistent scheduler)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=__REPO_DIR__/.env
ExecStart=__REPO_DIR__/.venv/bin/python -m agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
