# deploy/agent.service
[Unit]
Description=__AGENT_DISPLAY_NAME__ — Telegram ↔ Claude agent (long-poll bot + persistent scheduler)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=__REPO_DIR__/.env
ExecStart=__REPO_DIR__/.venv/bin/python -m agent
Restart=always
RestartSec=5
# A deploy/restart must not kill an in-flight claude turn: SIGTERM goes to the bot process
# only, PTB finishes the running handler (the turn completes and replies) before exiting.
# The subprocess otherwise dies with exit 143 and the chat gets "claude failed". Only after
# TimeoutStopSec does systemd SIGKILL the whole cgroup.
KillMode=mixed
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
